# sources/services/cwa-automated.service
[Unit]
Description=Calibre-Web Automated
After=network.target

[Service]
Type=simple
User=cwa
Group=cwa
WorkingDirectory=/opt/cwa-automated
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/cwa-automated/venv/bin/python /opt/cwa-automated/cps.py --host 0.0.0.0 --port 8083
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
